# systemd/backup-verify.service
[Unit]
Description=Verify backup integrity (light check)
After=network-online.target
Wants=network-online.target
OnFailure=backup-onfailure@%n.service

[Service]
Type=oneshot
ExecStart=/opt/miniforge/envs/backup/bin/backup verify
Nice=15
IOSchedulingClass=idle

# Bound a hung run so it can't hold a repo lock indefinitely (see backup.service).
TimeoutStartSec=1h

# Security hardening
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/backup
NoNewPrivileges=true
