# systemd/backup-verify-deep.timer
[Unit]
Description=Verify backup integrity weekly (deep check)

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3h
Persistent=true

[Install]
WantedBy=timers.target
